# data/systemd/ims.service.in
[Unit]
Description=Image Manipulation Service
AssertPathExists=@sysconfdir@

[Service]
ExecStart=@bindir@/ims --config @sysconfdir@/ims.conf

[Install]
WantedBy=default.target
